# ot-br-posix/files/etc/hotplug.d/net/30-wpan
#!/bin/sh

[ add = "$ACTION" ] && {
    [ "$INTERFACE" = "$(uci -q get otbr.bbr.ifname)" ] && {
        ubus send otbr-agent '{"id":9}'
    }
}
