# gl-sdk4-wireguard-ipv4-only/files/etc/uci-defaults/80-wg-noipv6
#!/bin/sh
# First-boot bootstrap: trigger an initial sync so existing routers
# reach the right state on first boot after the package is installed.

if [ -z "$IPKG_INSTROOT" ]; then
	/usr/sbin/wg-noipv6 sync >/dev/null 2>&1 || true
fi

exit 0
