# gl-sdk4-wireguard-ipv4-only/files/etc/init.d/wg-noipv6
#!/bin/sh /etc/rc.common
# Reload-only service: re-evaluate global IPv6 state on UCI changes
# to wireguard_server, network or firewall.

START=96
STOP=10
USE_PROCD=1

start_service() {
	# Hooks should only log to syslog -- never to stderr -- to keep procd
	# logs and serial consoles clean.
	WG_NOIPV6_QUIET=1 /usr/sbin/wg-noipv6 sync >/dev/null 2>&1 || true
}

reload_service() {
	WG_NOIPV6_QUIET=1 /usr/sbin/wg-noipv6 sync >/dev/null 2>&1 || true
}

stop_service() {
	# Routine `stop` may just be a network restart; full teardown
	# lives in `wg-noipv6 clear-all` (invoked by package prerm).
	return 0
}

service_triggers() {
	procd_add_reload_trigger "wireguard_server"
	procd_add_reload_trigger "network"
	procd_add_reload_trigger "firewall"
	procd_add_reload_trigger "glconfig"
}
