# gl-sdk4-wireguard-ipv4-only/files/etc/hotplug.d/net/99-wg-noipv6
#!/bin/sh
# The GL.iNet wireguard_server daemon brings its kernel device up
# without firing a netifd ifup, so the iface hotplug never sees it.
# This (netdev) hotplug is the only reliable hook in that case.

case "$ACTION" in
	add|register|"") ;;
	*) exit 0 ;;
esac

[ -n "$DEVICENAME" ] || exit 0
[ -d "/sys/class/net/$DEVICENAME" ] || exit 0

case "$DEVICENAME" in
	wg*|*wgserver*) ;;
	*) exit 0 ;;
esac

. /usr/lib/wg-noipv6/functions.sh

wg_iface_is_valid "$DEVICENAME" || exit 0

# Hotplug context: never echo to stderr, just syslog.
export WG_NOIPV6_QUIET=1

if [ "$(wg_global_ipv6_enabled)" = "0" ]; then
	/usr/sbin/wg-noipv6 apply "$DEVICENAME" >/dev/null 2>&1 || true
else
	# Global IPv6 was re-enabled but we still hold leftover host-side
	# enforcement or pending UCI snapshots. Run a full sync so the enable
	# path can clean up and restore.
	_stale=0
	[ -f "${WG_NOIPV6_SYSCTL_DIR}/99-wg-noipv6-${DEVICENAME}.conf" ] && _stale=1
	if [ "$_stale" -eq 0 ] && [ -d "$WG_NOIPV6_BACKUP_DIR" ]; then
		for _f in "$WG_NOIPV6_BACKUP_DIR"/*.snap; do
			[ -f "$_f" ] && { _stale=1; break; }
		done
	fi
	[ "$_stale" -eq 1 ] && /usr/sbin/wg-noipv6 sync >/dev/null 2>&1 || true
fi

exit 0
