# gl-sdk4-wireguard-ipv4-only/files/etc/hotplug.d/net/99-wg-noipv6
#!/bin/sh
# The GL.iNet wireguard_server daemon brings its kernel device up
# without firing a netifd ifup, so the iface hotplug never sees it.
# This (netdev) hotplug is the only reliable hook in that case.

case "$ACTION" in
	add|register|"") ;;
	*) exit 0 ;;
esac

[ -n "$DEVICENAME" ] || exit 0
[ -d "/sys/class/net/$DEVICENAME" ] || exit 0

case "$DEVICENAME" in
	wg*|*wgserver*) ;;
	*) exit 0 ;;
esac

. /usr/lib/wg-noipv6/functions.sh

wg_iface_is_valid "$DEVICENAME" || exit 0
[ "$(wg_global_ipv6_enabled)" = "0" ] || exit 0

/usr/sbin/wg-noipv6 apply "$DEVICENAME" >/dev/null 2>&1 || true

exit 0
