# ot-br-posix/files/etc/uci-defaults/13-ot-misc
#!/bin/sh

# Set log size
uci set system.@system[0].log_size=256
uci commit system

exit 0
